# configs/native/config.xml
<!--
    Sysmon Configuration - Schema 4.91
    Optimized for Wazuh 4.14.0
    Windows 11 Native Sysmon
-->

<Sysmon schemaversion="4.91">
  
  <HashAlgorithms>MD5,SHA256,IMPHASH</HashAlgorithms>
  <ArchiveDirectory>Sysmon</ArchiveDirectory>
  <CheckRevocation>False</CheckRevocation>
  <DnsLookup>True</DnsLookup>
  
  <EventFiltering>
    
    <!-- EVENT ID 1: Process Creation -->
    <RuleGroup name="ProcessCreate-Includes" groupRelation="or">
      <ProcessCreate onmatch="include">
        
        <!-- T1546.008 - Accessibility Features -->
        <ParentImage name="T1546.008" condition="image">sethc.exe</ParentImage>
        <ParentImage name="T1546.008" condition="image">utilman.exe</ParentImage>
        <ParentImage name="T1546.008" condition="image">osk.exe</ParentImage>
        
        <!-- T1548.002 - UAC Bypass -->
        <ParentImage name="T1548.002" condition="image">eventvwr.exe</ParentImage>
        <ParentImage name="T1548.002" condition="image">fodhelper.exe</ParentImage>
        
        <!-- T1003 - Credential Dumping -->
        <OriginalFileName name="T1003" condition="contains">procdump</OriginalFileName>
        <OriginalFileName name="T1003" condition="is">TTTracer.exe</OriginalFileName>
        <OriginalFileName name="T1003" condition="is">sqldumper.exe</OriginalFileName>
        <ParentImage name="T1003" condition="is">diskshadow.exe</ParentImage>
        
        <!-- T1059.001 - PowerShell -->
        <OriginalFileName name="T1059.001" condition="image">powershell.exe</OriginalFileName>
        <OriginalFileName name="T1059.001" condition="image">pwsh.exe</OriginalFileName>
        
        <!-- T1059.003 - CMD -->
        <OriginalFileName name="T1059.003" condition="is">cmd.exe</OriginalFileName>
        
        <!-- T1047 - WMI -->
        <OriginalFileName name="T1047" condition="is">wmic.exe</OriginalFileName>
        <ParentImage name="T1047" condition="image">wmiprvse.exe</ParentImage>
        
        <!-- T1053 - Scheduled Task -->
        <OriginalFileName name="T1053" condition="is">schtasks.exe</OriginalFileName>
        <OriginalFileName name="T1053" condition="is">at.exe</OriginalFileName>
        
        <!-- T1057 - Process Discovery -->
        <OriginalFileName name="T1057" condition="image">tasklist.exe</OriginalFileName>
        <OriginalFileName name="T1057" condition="image">qprocess.exe</OriginalFileName>
        
        <!-- T1105 - Ingress Tool Transfer -->
        <OriginalFileName name="T1105" condition="is">certutil.exe</OriginalFileName>
        <OriginalFileName name="T1105" condition="is">bitsadmin.exe</OriginalFileName>
        <OriginalFileName name="T1105" condition="is">curl.exe</OriginalFileName>
        
        <!-- T1218 - LOLBins -->
        <OriginalFileName name="T1218" condition="is">mshta.exe</OriginalFileName>
        <OriginalFileName name="T1218" condition="is">regsvr32.exe</OriginalFileName>
        <OriginalFileName name="T1218" condition="contains">rundll32.exe</OriginalFileName>
        <OriginalFileName name="T1218" condition="is">msiexec.exe</OriginalFileName>
        
        <!-- T1490 - Inhibit System Recovery -->
        <OriginalFileName name="T1490" condition="is">vssadmin.exe</OriginalFileName>
        <OriginalFileName name="T1490" condition="is">wbadmin.exe</OriginalFileName>
        <OriginalFileName name="T1490" condition="is">bcdedit.exe</OriginalFileName>
        
        <!-- T1036 - Masquerading -->
        <Image name="T1036" condition="begin with">C:\Temp\</Image>
        <Image name="T1036" condition="begin with">C:\Users\Public\</Image>
        <Image name="T1036" condition="begin with">C:\ProgramData\</Image>
        
        <!-- T1137 - Office Application Startup -->
        <ParentImage name="T1137" condition="image">winword.exe</ParentImage>
        <ParentImage name="T1137" condition="image">excel.exe</ParentImage>
        <ParentImage name="T1137" condition="image">powerpnt.exe</ParentImage>
        
        <!-- T1569.002 - Service Execution -->
        <OriginalFileName name="T1569.002" condition="contains">PsExec</OriginalFileName>
        
      </ProcessCreate>
    </RuleGroup>
    
    <RuleGroup name="ProcessCreate-Excludes" groupRelation="or">
      <ProcessCreate onmatch="exclude">
        <Image condition="is">C:\Windows\System32\svchost.exe</Image>
        <Image condition="is">C:\Windows\System32\conhost.exe</Image>
        <Image condition="begin with">C:\Program Files\Windows Defender\</Image>
      </ProcessCreate>
    </RuleGroup>
    
    
    <!-- EVENT ID 2: File Creation Time -->
    <RuleGroup name="FileCreateTime" groupRelation="or">
      <FileCreateTime onmatch="include">
        <TargetFilename name="T1070.006" condition="end with">.exe</TargetFilename>
        <Image name="T1070.006" condition="begin with">C:\Temp</Image>
      </FileCreateTime>
    </RuleGroup>
    
    
    <!-- EVENT ID 3: Network Connection -->
    <RuleGroup name="NetworkConnect-Includes" groupRelation="or">
      <NetworkConnect onmatch="include">
        
        <Image name="T1059.001" condition="image">powershell.exe</Image>
        <Image name="T1059.003" condition="image">cmd.exe</Image>
        <Image name="T1218" condition="image">mshta.exe</Image>
        <Image name="T1218" condition="image">regsvr32.exe</Image>
        <Image name="T1218" condition="image">rundll32.exe</Image>
        
        <DestinationPort name="T1571" condition="is">4444</DestinationPort>
        <DestinationPort name="T1571" condition="is">31337</DestinationPort>
        <DestinationPort name="T1021" condition="is">3389</DestinationPort>
        <DestinationPort name="T1021" condition="is">22</DestinationPort>
        
      </NetworkConnect>
    </RuleGroup>
    
    <RuleGroup name="NetworkConnect-Excludes" groupRelation="or">
      <NetworkConnect onmatch="exclude">
        <Image condition="is">C:\Windows\System32\svchost.exe</Image>
        <DestinationHostname condition="end with">.microsoft.com</DestinationHostname>
      </NetworkConnect>
    </RuleGroup>
    
    
    <!-- EVENT ID 5: Process Terminated -->
    <RuleGroup name="ProcessTerminate" groupRelation="or">
      <ProcessTerminate onmatch="include">
        <Image condition="begin with">C:\Users</Image>
        <Image condition="begin with">C:\Temp</Image>
      </ProcessTerminate>
    </RuleGroup>
    
    
    <!-- EVENT ID 6: Driver Loaded -->
    <RuleGroup name="DriverLoad" groupRelation="or">
      <DriverLoad onmatch="exclude">
        <Rule groupRelation="and">
          <Signature condition="contains">Microsoft</Signature>
          <SignatureStatus condition="is">Valid</SignatureStatus>
        </Rule>
      </DriverLoad>
    </RuleGroup>
    
    
    <!-- EVENT ID 7: Image Loaded -->
    <RuleGroup name="ImageLoad" groupRelation="or">
      <ImageLoad onmatch="include">
        <ImageLoaded name="T1055" condition="contains">vaultcli.dll</ImageLoaded>
        <ImageLoaded name="T1055" condition="contains">samlib.dll</ImageLoaded>
      </ImageLoad>
    </RuleGroup>
    
    
    <!-- EVENT ID 8: CreateRemoteThread -->
    <RuleGroup name="CreateRemoteThread-Includes" groupRelation="or">
      <CreateRemoteThread onmatch="include">
        <TargetImage name="T1055" condition="end with">.exe</TargetImage>
      </CreateRemoteThread>
    </RuleGroup>
    
    <RuleGroup name="CreateRemoteThread-Excludes" groupRelation="or">
      <CreateRemoteThread onmatch="exclude">
        <SourceImage condition="is">C:\Windows\System32\svchost.exe</SourceImage>
      </CreateRemoteThread>
    </RuleGroup>
    
    
    <!-- EVENT ID 9: RawAccessRead -->
    <RuleGroup name="RawAccessRead" groupRelation="or">
      <RawAccessRead onmatch="include">
        <Device name="T1006" condition="begin with">\\.\</Device>
      </RawAccessRead>
    </RuleGroup>
    
    
    <!-- EVENT ID 10: ProcessAccess -->
    <RuleGroup name="ProcessAccess-Includes" groupRelation="or">
      <ProcessAccess onmatch="include">
        <TargetImage name="T1003.001" condition="image">lsass.exe</TargetImage>
      </ProcessAccess>
    </RuleGroup>
    
    <RuleGroup name="ProcessAccess-Excludes" groupRelation="or">
      <ProcessAccess onmatch="exclude">
        <SourceImage condition="is">C:\Windows\System32\svchost.exe</SourceImage>
        <SourceImage condition="end with">Windows Defender\MsMpEng.exe</SourceImage>
      </ProcessAccess>
    </RuleGroup>
    
    
    <!-- EVENT ID 11: FileCreate -->
    <RuleGroup name="FileCreate" groupRelation="or">
      <FileCreate onmatch="include">
        
        <TargetFilename name="T1547.001" condition="contains">\Startup\</TargetFilename>
        <TargetFilename name="T1053" condition="begin with">C:\Windows\System32\Tasks\</TargetFilename>
        <TargetFilename name="T1036" condition="begin with">C:\Windows\Temp\</TargetFilename>
        <TargetFilename name="T1036" condition="begin with">C:\Users\Public\</TargetFilename>
        
        <TargetFilename name="T1059" condition="end with">.ps1</TargetFilename>
        <TargetFilename name="T1059" condition="end with">.vbs</TargetFilename>
        <TargetFilename name="T1059" condition="end with">.bat</TargetFilename>
        
      </FileCreate>
    </RuleGroup>
    
    
    <!-- EVENT ID 12/13/14: Registry Events -->
    <RuleGroup name="RegistryEvent" groupRelation="or">
      <RegistryEvent onmatch="include">
        
        <TargetObject name="T1547.001" condition="contains">\CurrentVersion\Run</TargetObject>
        <TargetObject name="T1543.003" condition="contains">\CurrentControlSet\Services\</TargetObject>
        <TargetObject name="T1562.001" condition="contains">\Windows Defender\</TargetObject>
        
      </RegistryEvent>
    </RuleGroup>
    
    
    <!-- EVENT ID 15: FileCreateStreamHash -->
    <RuleGroup name="FileCreateStreamHash-Includes" groupRelation="or">
      <FileCreateStreamHash onmatch="include">
        <TargetFilename name="T1564.004" condition="contains">:</TargetFilename>
      </FileCreateStreamHash>
    </RuleGroup>
    
    <RuleGroup name="FileCreateStreamHash-Excludes" groupRelation="or">
      <FileCreateStreamHash onmatch="exclude">
        <TargetFilename condition="end with">:Zone.Identifier</TargetFilename>
      </FileCreateStreamHash>
    </RuleGroup>
    
    
    <!-- EVENT ID 17/18: Pipe Events -->
    <RuleGroup name="PipeEvent" groupRelation="or">
      <PipeEvent onmatch="include">
        <PipeName name="CobaltStrike" condition="contains">msagent_</PipeName>
        <PipeName name="CobaltStrike" condition="contains">MSSE-</PipeName>
        <PipeName name="PSExec" condition="is">psexec</PipeName>
      </PipeEvent>
    </RuleGroup>
    
    
    <!-- EVENT ID 19/20/21: WMI Events -->
    <RuleGroup name="WmiEvent" groupRelation="or">
      <WmiEvent onmatch="include">
        <Operation name="T1546.003" condition="is">Created</Operation>
      </WmiEvent>
    </RuleGroup>
    
    
    <!-- EVENT ID 22: DNS Query -->
    <RuleGroup name="DnsQuery-Includes" groupRelation="or">
      <DnsQuery onmatch="include">
        
        <QueryName name="T1071.004" condition="end with">.tk</QueryName>
        <QueryName name="T1071.004" condition="end with">.ml</QueryName>
        <QueryName name="T1071.004" condition="end with">.ga</QueryName>
        
        <Image name="T1071.004" condition="image">powershell.exe</Image>
        <Image name="T1071.004" condition="image">cmd.exe</Image>
        
      </DnsQuery>
    </RuleGroup>
    
    <RuleGroup name="DnsQuery-Excludes" groupRelation="or">
      <DnsQuery onmatch="exclude">
        <QueryName condition="end with">microsoft.com</QueryName>
        <QueryName condition="end with">windows.com</QueryName>
        <Image condition="is">C:\Windows\System32\svchost.exe</Image>
      </DnsQuery>
    </RuleGroup>
    
    
    <!-- EVENT ID 23: FileDelete -->
    <RuleGroup name="FileDelete" groupRelation="or">
      <FileDelete onmatch="include">
        <TargetFilename name="T1070.004" condition="end with">.exe</TargetFilename>
        <TargetFilename name="T1070.004" condition="end with">.dll</TargetFilename>
        <TargetFilename name="T1070.004" condition="end with">.ps1</TargetFilename>
      </FileDelete>
    </RuleGroup>
    
    
    <!-- EVENT ID 24: ClipboardChange -->
    <RuleGroup name="ClipboardChange" groupRelation="or">
      <ClipboardChange onmatch="exclude">
        <Image condition="end with">chrome.exe</Image>
        <Image condition="image">WINWORD.EXE</Image>
        <Image condition="image">EXCEL.EXE</Image>
      </ClipboardChange>
    </RuleGroup>
    
    
    <!-- EVENT ID 25: ProcessTampering -->
    <RuleGroup name="ProcessTampering" groupRelation="or">
      <ProcessTampering onmatch="include">
        <Image name="T1055.012" condition="end with">.exe</Image>
      </ProcessTampering>
    </RuleGroup>
    
    
    <!-- EVENT ID 26: FileDeleteDetected -->
    <RuleGroup name="FileDeleteDetected" groupRelation="or">
      <FileDeleteDetected onmatch="include">
        <TargetFilename name="T1070.004" condition="end with">.exe</TargetFilename>
      </FileDeleteDetected>
    </RuleGroup>
    
    
    <!-- EVENT ID 29: FileExecutableDetected -->
    <RuleGroup name="FileExecutableDetected" groupRelation="or">
      <FileExecutableDetected onmatch="include">
        <TargetFilename name="T1204" condition="contains">\Downloads\</TargetFilename>
        <TargetFilename name="T1204" condition="begin with">C:\Users\Public\</TargetFilename>
      </FileExecutableDetected>
    </RuleGroup>
    
  </EventFiltering>
  
</Sysmon>
