# configs/native/sysmon-native-wazuh.xml
<!--
    ============================================================
    Sysmon Native Configuration - Windows 11 (Schema 4.91)
    ============================================================

    Project     : Unified Sysmon Configs
    Author      : m0us3r
    Repository  : https://github.com/mym0us3r/Unified-Sysmon-Configs
    Version     : 1.0
    Date        : 2026-03-23
    License     : GPL-3.0

    Platform    : Windows 11 24H2+ (Build 26100+) | Windows Server 2025
    Sysmon      : Native (Windows Optional Feature)
    Schema      : 4.91
    Integration : Wazuh SIEM/XDR 4.14.0+

    Description :
        Enterprise-grade Sysmon configuration for the Windows Native Sysmon
        feature. Tuned for high-fidelity threat detection with minimal noise,
        aligned with MITRE ATT&CK framework. Optimized for Wazuh ingestion
        and SOC/ThreatHunting/CSIRT workflows.

    Covered Event IDs :
        EID  1  - Process Creation         (T1059, T1218, T1003, T1036)
        EID  2  - File Creation Time       (T1070.006)
        EID  3  - Network Connection       (T1571, T1021, T1059)
        EID  5  - Process Terminated       (User/Temp paths)
        EID  6  - Driver Loaded            (Exclude valid Microsoft/signed)
        EID  7  - Image Loaded             (T1055 - DLL Injection)
        EID  8  - CreateRemoteThread       (T1055 - Process Injection)
        EID  9  - RawAccessRead            (T1006)
        EID 10  - ProcessAccess            (T1003.001 - LSASS)
        EID 11  - FileCreate               (T1547, T1053, T1036, T1059)
        EID 12/13/14 - RegistryEvent       (T1547.001, T1543.003, T1562.001)
        EID 15  - FileCreateStreamHash     (T1564.004 - ADS)
        EID 17/18 - PipeEvent              (CobaltStrike, PSExec)
        EID 19/20/21 - WmiEvent            (T1546.003)
        EID 22  - DnsQuery                 (T1071.004 - C2 over DNS)
        EID 23  - FileDelete               (T1070.004)
        EID 24  - ClipboardChange          (Exclude browsers/Office)
        EID 25  - ProcessTampering         (T1055.012 - Hollowing)
        EID 26  - FileDeleteDetected       (T1070.004)
        EID 29  - FileExecutableDetected   (T1204 - Downloads/Public)

    Notes :
        - This configuration uses the RuleGroup nested syntax (Schema 4.x).
          EventFiltering > RuleGroup > EventType onmatch="include|exclude"
        - Sysmon Native does NOT coexist with standalone Sysmon.
          Remove the standalone version before enabling the feature.
        - Apply with: sysmon -c sysmon-windows-native-wazuh.xml
        - Validate with: Check-SysmonHealth.ps1 (included in /scripts)
        - Config hash verified after each deployment via EID 16.
        - Wazuh rule alignment: /ruleset/rules/wazuh-server-4.14/

    References :
        - https://learn.microsoft.com/en-us/windows/security/operating-system-security/sysmon/how-to-enable-sysmon
        - https://techcommunity.microsoft.com/blog/windows-itpro-blog/native-sysmon-functionality-coming-to-windows/4468112
        - https://github.com/SwiftOnSecurity/sysmon-config
        - https://github.com/olafhartong/sysmon-modular
        - https://attack.mitre.org/

    ============================================================
-->

<Sysmon schemaversion="4.91">
  
  <HashAlgorithms>MD5,SHA256,IMPHASH</HashAlgorithms>
  <ArchiveDirectory>Sysmon</ArchiveDirectory>
  <CheckRevocation>False</CheckRevocation>
  <DnsLookup>True</DnsLookup>

  <EventFiltering>
    <RuleGroup name="ProcessCreate-Includes" groupRelation="or">
      <ProcessCreate onmatch="include">
        <Rule name="T1546.008-AccessibilityFeatures" groupRelation="and">
          <Image condition="contains any">sethc.exe;utilman.exe;osk.exe;Magnify.exe;DisplaySwitch.exe;AtBroker.exe</Image>
          <ParentImage condition="end with">cmd.exe</ParentImage>
        </Rule>
        <Rule name="T1548.002-UACBypass" groupRelation="and">
          <ParentImage condition="contains any">eventvwr.exe;fodhelper.exe;computerdefaults.exe;sdclt.exe</ParentImage>
          <Image condition="is">cmd.exe</Image>
        </Rule>
        <Rule name="T1003-CredentialDumping" groupRelation="and">
          <OriginalFileName condition="contains any">procdump;TTTracer.exe;sqldumper.exe;PypyKatz;Dumpert</OriginalFileName>
          <CommandLine condition="contains any">lsass;1.dmp;full;mini;all;accept-eula</CommandLine>
        </Rule>
        <Rule name="T1059.001-PowerShell-Execution" groupRelation="and">
          <OriginalFileName condition="contains any">powershell.exe;pwsh.exe;powershell_ise.exe</OriginalFileName>
          <CommandLine condition="contains any">-enc;-nop;-w hidden;IEX;Invoke-Expression;DownloadString;FromBase64String;bypass;-ExecutionPolicy</CommandLine>
        </Rule>
        <Rule name="T1059.003-CMD-Execution" groupRelation="and">
          <OriginalFileName condition="is">cmd.exe</OriginalFileName>
          <CommandLine condition="contains any">/c;/k;certutil;bitsadmin;powershell;net user;net localgroup;whoami</CommandLine>
        </Rule>
        <Rule name="T1047-WMI-Abuse" groupRelation="and">
          <OriginalFileName condition="is">wmic.exe</OriginalFileName>
          <CommandLine condition="contains any">process call create;shadowcopy delete;node;context;useraccount</CommandLine>
        </Rule>
        <Rule name="T1053-ScheduledTask" groupRelation="and">
          <OriginalFileName condition="contains any">schtasks.exe;at.exe</OriginalFileName>
          <CommandLine condition="contains any">/create;/delete;/change;/run;/xml</CommandLine>
        </Rule>
        <Rule name="T1105-IngressToolTransfer" groupRelation="and">
          <OriginalFileName condition="contains any">certutil.exe;bitsadmin.exe;curl.exe;wget.exe</OriginalFileName>
          <CommandLine condition="contains any">-urlcache;-split;-decode;/transfer;/addfile;-o;-O;--output</CommandLine>
        </Rule>
        <Rule name="T1218-LOLBins" groupRelation="and">
          <OriginalFileName condition="contains any">mshta.exe;regsvr32.exe;rundll32.exe;msiexec.exe;control.exe</OriginalFileName>
          <CommandLine condition="contains any">javascript;url.dll;/i http;.hta;scrobj.dll;pcwrun.exe</CommandLine>
        </Rule>
        <Rule name="T1490-InhibitRecovery" groupRelation="and">
          <OriginalFileName condition="contains any">vssadmin.exe;wbadmin.exe;bcdedit.exe</OriginalFileName>
          <CommandLine condition="contains any">delete shadows;delete systemstatebackup;recoveryenabled No;ignoreallfailures</CommandLine>
        </Rule>
        <Rule name="T1036-Masquerading" groupRelation="and">
          <Image condition="begin with">C:\Temp\;C:\Users\Public\;C:\ProgramData\;C:\PerfLogs\</Image>
          <Image condition="end with">.exe;.dll;.bat;.ps1;.vbs;.js</Image>
        </Rule>
        <Rule name="T1137-OfficeStartup" groupRelation="and">
          <ParentImage condition="contains any">winword.exe;excel.exe;powerpnt.exe;outlook.exe</ParentImage>
          <Image condition="contains any">cmd.exe;powershell.exe;rundll32.exe;mshta.exe;wscript.exe;cscript.exe</Image>
        </Rule>
        <Rule name="T1569.002-ServiceExecution" groupRelation="and">
          <OriginalFileName condition="contains">PsExec</OriginalFileName>
          <CommandLine condition="contains any">-s;-x;-accepteula;-i;-d</CommandLine>
        </Rule>
      </ProcessCreate>
    </RuleGroup>

    <RuleGroup name="ProcessCreate-Excludes" groupRelation="or">
      <ProcessCreate onmatch="exclude">
        <Image condition="is any">C:\Windows\System32\svchost.exe;C:\Windows\System32\conhost.exe;C:\Windows\System32\taskhostw.exe</Image>
        <Image condition="begin with">C:\Program Files\Windows Defender\;C:\Program Files\Wazuh-Agent\</Image>
      </ProcessCreate>
    </RuleGroup>

    <RuleGroup name="FileCreateTime" groupRelation="or">
      <FileCreateTime onmatch="include">
        <Rule name="T1070.006-Timestomping" groupRelation="and">
          <TargetFilename condition="end with">.exe;.dll</TargetFilename>
          <Image condition="begin with">C:\Temp\;C:\Users\Public\</Image>
        </Rule>
      </FileCreateTime>
    </RuleGroup>

    <RuleGroup name="NetworkConnect-Includes" groupRelation="or">
      <NetworkConnect onmatch="include">
        <Rule name="T1059-SuspiciousOutbound" groupRelation="and">
          <Image condition="contains any">powershell.exe;cmd.exe;mshta.exe;regsvr32.exe;rundll32.exe;wscript.exe</Image>
          <DestinationPort condition="excludes any">80;443;53</DestinationPort>
        </Rule>
        <Rule name="T1571-C2-Ports" groupRelation="and">
          <Initiated condition="is">true</Initiated>
          <DestinationPort condition="is any">4444;31337;8080;8888;9001</DestinationPort>
        </Rule>
        <Rule name="T1021-LateralMovement" groupRelation="and">
          <Image condition="is">System</Image>
          <DestinationPort condition="is any">445;3389;22;5985;5986</DestinationPort>
        </Rule>
      </NetworkConnect>
    </RuleGroup>

    <RuleGroup name="NetworkConnect-Excludes" groupRelation="or">
      <NetworkConnect onmatch="exclude">
        <Image condition="is any">C:\Windows\System32\svchost.exe;C:\Program Files\Google\Chrome\Application\chrome.exe</Image>
        <DestinationHostname condition="end with">.microsoft.com;.windows.com;.delivery.mp.microsoft.com</DestinationHostname>
      </NetworkConnect>
    </RuleGroup>

    <RuleGroup name="ProcessTerminate" groupRelation="or">
      <ProcessTerminate onmatch="include">
        <Image condition="begin with">C:\Users\Public\;C:\Windows\Temp\;C:\Temp\</Image>
      </ProcessTerminate>
    </RuleGroup>

    <RuleGroup name="DriverLoad" groupRelation="or">
      <DriverLoad onmatch="exclude">
        <Rule groupRelation="and">
          <Signature condition="contains">Microsoft</Signature>
          <SignatureStatus condition="is">Valid</SignatureStatus>
        </Rule>
      </DriverLoad>
    </RuleGroup>

    <RuleGroup name="ImageLoad" groupRelation="or">
      <ImageLoad onmatch="include">
        <Rule name="T1055-CredentialAccess-DLL" groupRelation="and">
          <ImageLoaded condition="contains any">vaultcli.dll;samlib.dll;winscard.dll;cryptdll.dll</ImageLoaded>
          <Image condition="excludes any">C:\Windows\System32\lsass.exe;C:\Windows\System32\svchost.exe</Image>
        </Rule>
      </ImageLoad>
    </RuleGroup>

    <RuleGroup name="CreateRemoteThread-Includes" groupRelation="or">
      <CreateRemoteThread onmatch="include">
        <Rule name="T1055-Injection" groupRelation="and">
          <TargetImage condition="end with">.exe</TargetImage>
          <SourceImage condition="begin with">C:\Temp\;C:\Users\Public\;C:\PerfLogs\</SourceImage>
        </Rule>
      </CreateRemoteThread>
    </RuleGroup>

    <RuleGroup name="RawAccessRead" groupRelation="or">
      <RawAccessRead onmatch="include">
        <Rule name="T1006-DirectVolumeAccess" groupRelation="and">
          <Device condition="begin with">\\.\</Device>
          <Image condition="is not">C:\Windows\System32\svchost.exe</Image>
        </Rule>
      </RawAccessRead>
    </RuleGroup>

    <RuleGroup name="ProcessAccess-Includes" groupRelation="or">
      <ProcessAccess onmatch="include">
        <Rule name="T1003.001-LSASS-Access" groupRelation="and">
          <TargetImage condition="image">lsass.exe</TargetImage>
          <CallTrace condition="contains">unknown</CallTrace>
        </Rule>
      </ProcessAccess>
    </RuleGroup>

    <RuleGroup name="ProcessAccess-Excludes" groupRelation="or">
      <ProcessAccess onmatch="exclude">
        <SourceImage condition="contains any">C:\Windows\System32\svchost.exe;Windows Defender\MsMpEng.exe</SourceImage>
      </ProcessAccess>
    </RuleGroup>

    <RuleGroup name="FileCreate" groupRelation="or">
      <FileCreate onmatch="include">
        <Rule name="T1547.001-Persistence-Startup" groupRelation="and">
          <TargetFilename condition="contains">\Startup\</TargetFilename>
          <Image condition="is not">C:\Windows\explorer.exe</Image>
        </Rule>
        <Rule name="T1053-Persistence-Tasks" groupRelation="and">
          <TargetFilename condition="begin with">C:\Windows\System32\Tasks\</TargetFilename>
          <Image condition="is not">C:\Windows\System32\svchost.exe</Image>
        </Rule>
        <Rule name="T1059-SuspiciousScripts" groupRelation="and">
          <TargetFilename condition="end with">.ps1;.vbs;.bat;.js;.hta</TargetFilename>
          <TargetFilename condition="begin with">C:\Windows\Temp\;C:\Users\Public\</TargetFilename>
        </Rule>
      </FileCreate>
    </RuleGroup>

    <RuleGroup name="RegistryEvent" groupRelation="or">
      <RegistryEvent onmatch="include">
        <Rule name="T1562.001-DisableSecurity" groupRelation="and">
          <TargetObject condition="contains any">\Windows Defender\;\Reporting\DisableAntiSpyware</TargetObject>
          <Details condition="contains any">DisableAntiSpyware;DisableRealtimeMonitoring;DisableBehaviorMonitoring</Details>
        </Rule>
        <Rule name="T1547.001-RunKey-Persistence" groupRelation="and">
          <TargetObject condition="contains any">\CurrentVersion\Run;\CurrentVersion\RunOnce</TargetObject>
          <Image condition="excludes any">C:\Windows\explorer.exe;C:\Windows\System32\msiexec.exe</Image>
        </Rule>
      </RegistryEvent>
    </RuleGroup>

    <RuleGroup name="FileCreateStreamHash-Includes" groupRelation="or">
      <FileCreateStreamHash onmatch="include">
        <Rule name="T1564.004-ADS-Creation" groupRelation="and">
          <TargetFilename condition="contains">:</TargetFilename>
          <Image condition="begin with">C:\Temp\;C:\Users\Public\</Image>
        </Rule>
      </FileCreateStreamHash>
    </RuleGroup>

    <RuleGroup name="FileCreateStreamHash-Excludes" groupRelation="or">
      <FileCreateStreamHash onmatch="exclude">
        <TargetFilename condition="end with">:Zone.Identifier</TargetFilename>
      </FileCreateStreamHash>
    </RuleGroup>

    <RuleGroup name="PipeEvent" groupRelation="or">
      <PipeEvent onmatch="include">
        <Rule name="C2-Known-Pipes" groupRelation="and">
          <PipeName condition="contains any">msagent_;MSSE-;psexec;postex_;status_</PipeName>
        </Rule>
      </PipeEvent>
    </RuleGroup>

    <RuleGroup name="WmiEvent" groupRelation="or">
      <WmiEvent onmatch="include">
        <Rule name="T1546.003-WMI-Persistence" groupRelation="and">
          <Operation condition="is">Created</Operation>
          <EventNamespace condition="contains">subscription</EventNamespace>
        </Rule>
      </WmiEvent>
    </RuleGroup>

    <RuleGroup name="DnsQuery-Includes" groupRelation="or">
      <DnsQuery onmatch="include">
        <Rule name="T1071.004-Suspicious-TLDs" groupRelation="and">
          <QueryName condition="end with">.tk;.ml;.ga;.cf;.gq</QueryName>
          <Image condition="contains any">powershell.exe;cmd.exe;wscript.exe;cscript.exe</Image>
        </Rule>
      </DnsQuery>
    </RuleGroup>

    <RuleGroup name="DnsQuery-Excludes" groupRelation="or">
      <DnsQuery onmatch="exclude">
        <QueryName condition="end with">microsoft.com;windows.com;akadns.net;azure.com</QueryName>
        <Image condition="is">C:\Windows\System32\svchost.exe</Image>
      </DnsQuery>
    </RuleGroup>

    <RuleGroup name="FileDelete" groupRelation="or">
      <FileDelete onmatch="include">
        <Rule name="T1070.004-FileDeletion" groupRelation="and">
          <TargetFilename condition="end with">.exe;.dll;.ps1</TargetFilename>
          <Image condition="begin with">C:\Temp\;C:\Users\Public\</Image>
        </Rule>
      </FileDelete>
    </RuleGroup>

    <RuleGroup name="ClipboardChange" groupRelation="or">
      <ClipboardChange onmatch="exclude">
        <Image condition="end with">chrome.exe;msedge.exe;firefox.exe</Image>
        <Image condition="contains any">WINWORD.EXE;EXCEL.EXE;POWERPNT.EXE</Image>
      </ClipboardChange>
    </RuleGroup>

    <RuleGroup name="ProcessTampering" groupRelation="or">
      <ProcessTampering onmatch="include">
        <Rule name="T1055.012-Hollowing" groupRelation="and">
          <Image condition="end with">.exe</Image>
          <Type condition="is">Image replaced</Type>
        </Rule>
      </ProcessTampering>
    </RuleGroup>

    <RuleGroup name="FileDeleteDetected" groupRelation="or">
      <FileDeleteDetected onmatch="include">
        <Rule name="T1070.004-Wipe-Detection" groupRelation="and">
          <TargetFilename condition="end with">.exe</TargetFilename>
          <Image condition="begin with">C:\Temp\;C:\Users\Public\</Image>
        </Rule>
      </FileDeleteDetected>
    </RuleGroup>

    <RuleGroup name="FileExecutableDetected" groupRelation="or">
      <FileExecutableDetected onmatch="include">
        <Rule name="T1204-New-Executable" groupRelation="and">
          <TargetFilename condition="contains any">\Downloads\;C:\Users\Public\</TargetFilename>
          <Hashes condition="contains">IMPHASH</Hashes>
        </Rule>
      </FileExecutableDetected>
    </RuleGroup>
  </EventFiltering>
</Sysmon>
